# deploy/webhook.template
#!/sbin/openrc-run

command=webhook
command_args="-hooks $PROJECT_ROOT/deploy/hooks.json &"